# operations/automation/launchd/com.aiagent.tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.aiagent.tunnel</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>/Users/gedizucar/Projects/ai-agent-team-mobile-app/operations/automation/scripts/start_tunnel.sh</string>
  </array>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/aiagent-tunnel.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/aiagent-tunnel.err.log</string>
</dict>
</plist>
